# containers/docker/runit/dovecot/run
#!/bin/bash

exec /usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
